# drop-in/drop-in-test.service
[Unit]
Description = Drop in test

[Service]
Type = simple
ExecStart = /usr/local/bin/drop-in.sh $env1 $env2

[Install]
WantedBy = multi-user.target
